# budget-server.conf
description "budget-web-server"

start on started networking
stop on runlevel [!2345]

exec su -u budget -c "/home/myusername/run_at_startup.sh"

